# infra/hermes-agents/systemd/hermes-gateway-sakjules.service
[Unit]
Description=Hermes Agent Gateway - SakJules
After=network.target

[Service]
# Replace 'beer' with the Linux user that will run the agent.
User=beer
Group=beer

# The full path to the Hermes executable.
# This assumes Hermes is installed in the user's local bin.
ExecStart=/home/beer/.local/bin/hermes --profile sakjules

# Set the SAKTHAI_HOME environment variable so the agent can find the shared memory database.
Environment="SAKTHAI_HOME=/home/beer/.sakthai"
Environment="SAKTHAI_SHELL_ALLOW=1"

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target