# infra/hermes-agents/systemd/hermes-gateway-saksee.service.d/display.conf
# Allow Saksee (backend: local) to open headed GUI browsers on the WSLg display.
[Service]
Environment="DISPLAY=:0"
Environment="WAYLAND_DISPLAY=wayland-0"
